# systemd/etc/systemd/system/set_swappines.service
[Unit]
Description=Set Swappiness

[Service]
ExecStart=/usr/bin/set_swappines.sh

[Install]
WantedBy=multi-user.target

