# systemd/etc/systemd/system/swap_opt_cmd.service
[Unit]
Description=Swap Opt Command Key

[Service]
ExecStart=/usr/bin/swap_opt_cmd.sh

[Install]
WantedBy=multi-user.target
